# engines/moss-transcribe/com.voxstudio.moss-transcribe.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.voxstudio.moss-transcribe</string>

  <key>ProgramArguments</key>
  <array>
    <string>@SERVER_DIR@/.venv/bin/python</string>
    <string>-m</string>
    <string>uvicorn</string>
    <string>server_moss:app</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>18087</string>
  </array>

  <key>WorkingDirectory</key>
  <string>@SERVER_DIR@</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>MOSS_LIBRARY</key>
    <string>@MOSS_LIBRARY@</string>
    <key>MOSS_MODEL</key>
    <string>@MOSS_MODEL@</string>
    <key>MOSS_MODEL_NAME</key>
    <string>moss-transcribe-diarize</string>
    <key>MOSS_MAX_UPLOAD_BYTES</key>
    <string>1073741824</string>
    <key>MOSS_MAX_DURATION_SECONDS</key>
    <string>7200</string>
    <key>MOSS_QUEUE_LIMIT</key>
    <string>4</string>
    <key>MOSS_FFMPEG_TIMEOUT_SECONDS</key>
    <string>300</string>
    <key>MOSS_FFMPEG</key>
    <string>@MOSS_FFMPEG@</string>
    <key>MOSS_FFPROBE</key>
    <string>@MOSS_FFPROBE@</string>
    <key>GGML_METAL_NO_RESIDENCY</key>
    <string>1</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProcessType</key>
  <string>Interactive</string>

  <key>StandardOutPath</key>
  <string>@LOG_DIR@/moss-transcribe.log</string>
  <key>StandardErrorPath</key>
  <string>@LOG_DIR@/moss-transcribe.error.log</string>
</dict>
</plist>
