# engines/llamacpp/com.voxstudio.gemma-llm.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- launchd user agent for the local conversation LLM (llama-server + Gemma QAT gguf).
     Replace %LLAMA_SERVER% (e.g. /opt/homebrew/bin/llama-server) and %MODEL_GGUF%. -->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.voxstudio.gemma-llm</string>
  <key>ProgramArguments</key>
  <array>
    <string>%LLAMA_SERVER%</string>
    <string>-m</string><string>%MODEL_GGUF%</string>
    <string>--host</string><string>127.0.0.1</string>
    <string>--port</string><string>18094</string>
    <string>-ngl</string><string>99</string>
    <string>-c</string><string>16384</string>
    <string>-np</string><string>1</string>
    <string>--jinja</string>
    <string>--flash-attn</string><string>auto</string>
    <string>--reasoning-budget</string><string>0</string>
    <string>-a</string><string>gemma4-12b-qat</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/tmp/gemma-llm.log</string>
  <key>StandardErrorPath</key><string>/tmp/gemma-llm.log</string>
</dict>
</plist>
