# engines/voxcpm2-server/voxcpm2-tts.service
[Unit]
Description=VoxCPM2 TTS HTTP service (voxstudio engine)
After=network-online.target

[Service]
# Adjust WorkingDirectory and VOXCPM2_* env to your host.
WorkingDirectory=%h/voxstudio/engines/voxcpm2-server
Environment=VOXCPM2_BASE=%h/tts-eval-voxcpm2
# Set this to the deployment-generated manifest SHA-256 to make design profiles
# reproducible across model artifact changes.
# Environment=VOXCPM2_MODEL_MANIFEST_SHA256=<64 lowercase hex characters>
Environment=PYTHONUNBUFFERED=1
# The caching allocator's reserved pool ratchets up with the length of one generation and
# never comes back down, which a co-tenant model on the same card pays for. Expandable
# segments map and unmap on demand instead of reserving fixed-size blocks, cutting the pool
# a long generation needs by ~10x. Must be set before torch initialises CUDA, which is why
# it lives here and not in the code. Pairs with the `empty_cache()` in `_generate`.
Environment=PYTORCH_CUDA_ALLOC_CONF=expandable_segments:True
# Run the interpreter that `uv sync --frozen` installed, directly. Do NOT wrap this
# in `uv run`: under systemd, uv swallows the child's stdout and stderr, so uvicorn's
# access log and any traceback never reach the journal. uv installs; it does not launch.
# Path is wherever UV_PROJECT_ENVIRONMENT pointed during sync (default: ./.venv).
ExecStart=%h/voxstudio/engines/voxcpm2-server/.venv/bin/python -m uvicorn server_voxcpm2:app --host 0.0.0.0 --port 8880
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
